# remove-download/remove-downloads.timer
[Unit]
Description=Run remove-downloads after boot

[Timer]
OnBootSec=15min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
