# remove-downloads/remove-downloads.service.template
[Unit]
Description=Remove contents of ~/Downloads

[Service]
Type=oneshot
ExecStart={{PATH}}/remove-downloads.sh

[Install]
WantedBy=multi-user.target
