# remove-downloads/remove-downloads.timer
[Unit]
Description=Run remove-downloads after boot

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
